# RemoteScheduling/03 Start Stream.xml
<?xml version="1.0" encoding="UTF-16"?> 
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task"> 
  <RegistrationInfo> 
    <Date>2023-04-02T09:44:00.0000</Date> 
    <Author>W10-D-24\User</Author>
    <URI>04 08 2023\03 Start Stream</URI>
  </RegistrationInfo>
  <Triggers>
    <TimeTrigger>
      <StartBoundary>2023-04-08T16:59:00</StartBoundary>
      <Enabled>true</Enabled>
    </TimeTrigger>
  </Triggers>
  <Actions Context="Author"> 
    <Exec> 
      <Command>C:\Users\User\Desktop\startStream.bat</Command> 
      <Arguments></Arguments> 
      <WorkingDirectory>C:\Users\User\Desktop</WorkingDirectory> 
    </Exec>
  </Actions> 
</Task> 